# ndash.service
[Unit]
Description=NDash PowerDNS Admin Dashboard
Documentation=https://github.com/yourusername/ndash
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/ndash
Environment=NODE_ENV=production
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=ndash

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/ndash
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
